# SPECS.r/rsh/rexec.socket
[Unit]
Description=Remote Execution Facilities Activation Socket

[Socket]
ListenStream=512
Accept=true

[Install]
WantedBy=sockets.target
