# SPECS.v/vsftpd/vsftpd.target
[Unit]
Description=FTP daemon
After=network.target

[Install]
WantedBy=multi-user.target
