# SPECS.i/inn/innd-expire.timer
[Unit]
Description=Inetnet News Daemon Exp

[Timer]
OnBootSec=3min
OnUnitActiveSec=1d
Unit=innd-expire.service

[Install]
WantedBy=multi-user.target