# SPECS.o/openldap/slapd.service
[Unit]
Description=OpenLDAP Server Daemon
After=syslog.target network.target
Documentation=man:slapd
Documentation=man:slapd-config
Documentation=man:slapd-hdb
Documentation=man:slapd-mdb
Documentation=file:///usr/share/doc/openldap-servers/guide.html

[Service]
Type=forking
PIDFile=/var/run/openldap/slapd.pid
ExecStartPre=/usr/libexec/openldap/check-config.sh
ExecStart=/usr/sbin/slapd -u ldap -h "ldap:/// ldapi:///"

[Install]
WantedBy=multi-user.target
Alias=openldap.service
