# SPECS.c/clamav/clamd.SERVICE.init
#!/bin/bash
#
# chkconfig: - 75 35
# description: The clamd server running for <SERVICE>

CLAMD_SERVICE=<SERVICE>
. /usr/share/clamav/clamd-wrapper
