# SPECS.g/geoipupdate/geoipupdate.cron
#!/bin/sh

# This updates the IPv4 databases based on product numbers in /etc/GeoIP.conf
geoipupdate > /dev/null
